# artifacts/rootfs/overlay/etc/init.d/guest-agent
#!/sbin/openrc-run

description="AgentBox Guest Agent"
command="/usr/local/bin/guest-agent"
command_background=true
pidfile="/run/guest-agent.pid"
output_log="/var/log/guest-agent.log"
error_log="/var/log/guest-agent.log"

depend() {
    after localmount
}

start_pre() {
    # Bring up loopback so port forwarding (vsock -> 127.0.0.1) works
    # even when full networking is disabled
    /sbin/ip link set lo up 2>/dev/null || true
}
